# services/aploader/xinetd.d/aploader
service aploader
{
	port		= 4242
	id			= aploader
	server		= /home/aploader/run.sh
	socket_type	= stream
	type		= UNLISTED
	protocol	= tcp
	user		= aploader
	wait		= no
}
