# etc/systemd/system/pfs.service
[Unit]
Description=PostgreSQL File System
After=network.target

[Service]
Type=simple
User=pfs
Group=pfs

ExecStart=/home/pfs.rb /home/pfs/t
